# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log
priority=10

[program:backend]
command=node /app/server/dist/index.js
directory=/app/server
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
environment=NODE_ENV="production",DATABASE_URL="postgresql://milicard:%(ENV_DB_PASSWORD)s@localhost:5432/milicard?schema=public",JWT_SECRET="%(ENV_JWT_SECRET)s",JWT_REFRESH_SECRET="%(ENV_JWT_SECRET)s_refresh",PORT="3000",CORS_ORIGIN="*"
priority=20
startsecs=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=30
